# dehydrated.crontab
* * * * * root /usr/bin/dehydrated --cron 1>/dev/stdout 2>/dev/stderr
